# etc/supervisor/conf.d/my-flask-app.conf
[program:my-flask-app]
user=www-data
autostart=true
autorestart=true
command=/opt/my-flask-app/venv/bin/uwsgi /opt/my-flask-app/uwsgi.yaml